# infra/systemd/friendgroup-api.service
[Unit]
Description=GEM API (Fastify dev server)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/gem
ExecStart=/usr/bin/npm run dev:api
Restart=always
RestartSec=3
Environment=NODE_ENV=development

[Install]
WantedBy=default.target