# infra/systemd/friendgroup-web.service
[Unit]
Description=GEM Web (Vite dev server)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/gem
ExecStart=/usr/bin/npm run dev:web
Restart=always
RestartSec=3
Environment=NODE_ENV=development

[Install]
WantedBy=default.target