# infra/systemd/gem-web.prod.service
[Unit]
Description=Gem Web (production preview)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aidan
Group=devs
WorkingDirectory=/var/www/gem
Environment=NODE_ENV=production
Environment=VITE_API_BASE_URL=https://api-gem.aidanlenahan.com
ExecStart=/usr/bin/npm --workspace apps/web run preview -- --host 127.0.0.1 --port 4173 --strictPort
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/gem /tmp
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
