# infra/systemd/gem-api.prod.service
[Unit]
Description=Gem API (production)
After=network-online.target docker.service
Wants=network-online.target docker.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=aidan
Group=devs
WorkingDirectory=/var/www/gem
EnvironmentFile=/etc/gem/gem-api.env
ExecStartPre=/usr/bin/mkdir -p /var/www/gem/apps/api/dist/generated
ExecStartPre=/usr/bin/cp -r /var/www/gem/apps/api/src/generated/prisma /var/www/gem/apps/api/dist/generated/
ExecStart=/usr/bin/node /var/www/gem/apps/api/dist/server.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/gem /tmp
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
